# deploy/greatwalk-bot.service
# GreatWalkBot system service (systemd)
#
# Production layout (see docs/deployment.md):
#   User:     greatwalk
#   Home:     /opt/greatwalk
#   Repo:     /opt/greatwalk-bot
#   Secrets:  /etc/greatwalk-bot/greatwalk-bot.env
#
# Install (after substituting placeholders):
#   sudo mkdir -p /etc/greatwalk-bot
#   sed ... deploy/greatwalk-bot.service > /tmp/greatwalk-bot.service
#   sudo cp /tmp/greatwalk-bot.service /etc/systemd/system/greatwalk-bot.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now greatwalk-bot
#
# Placeholders to replace before install:
#   __GREATWALKBOT_USER__     dedicated service account (e.g. greatwalk)
#   __GREATWALKBOT_HOME__      service account home (e.g. /opt/greatwalk)
#   __GREATWALKBOT_REPO__      application tree (e.g. /opt/greatwalk-bot)
#
# Runtime state under __GREATWALKBOT_REPO__:
#   config.yaml
#   logs/greatwalkbot.log
#   logs/status.json
#   data/seen.db

[Unit]
Description=GreatWalkBot availability watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__GREATWALKBOT_USER__
Group=__GREATWALKBOT_USER__
WorkingDirectory=__GREATWALKBOT_REPO__
ExecStart=__GREATWALKBOT_REPO__/.venv/bin/gwbot watch __GREATWALKBOT_REPO__/config.yaml
EnvironmentFile=-/etc/greatwalk-bot/greatwalk-bot.env
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
